# .config/systemd/user/offlineimap-oneshot.timer
[Unit]
Description=Run offlineimap every 30 minutes

[Timer]
OnBootSec=30min
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
